# playbooks/templates/supervisor-service.conf.j2
[program:{{name}}]
command={{command}}
autostart=true
autorestart=true
stderr_logfile=/var/log/{{name}}.err.log
stdout_logfile=/var/log/{{name}}.out.log
